# arkivbot.crontab
MAILTO=""
0 1 * * * jsub -N ArkivBot /data/project/arkivbot/run.job
