# arkivbot.crontab
MAILTO=""
0 1 * * * jsub /data/project/arkivbot/run.job
